# docker/etc/profile.d/common.sh
alias ll='ls -al --color=auto'

export CLICOLOR=1
export LSCOLORS=gxfxcxdxbxegedabagacad
export LANG=ja_JP.UTF-8
export LANGUAGE=ja_JP:ja
export LC_ALL=ja_JP.UTF-8
export SHELL=/usr/bin/zsh
export PS1="[\u@$CONTAINER_NAME \w]$"
